# deploy/researchbuilder.service
[Unit]
Description=ResearchBuilder API
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/researchbuilder/api
ExecStart=/opt/researchbuilder/.venv/bin/uvicorn index:app --host 127.0.0.1 --port 8000 --workers 1 --timeout-keep-alive 300
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment file
EnvironmentFile=/opt/researchbuilder/.env

# Security hardening
NoNewPrivileges=true
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/opt/researchbuilder/api/output /tmp

[Install]
WantedBy=multi-user.target
